# timer.cron
* * * * * php /Users/mustafa/Dropbox/Work/Criexe/vHosts/Base/timer.cx